# projects/ROCKNIX/packages/rocknix/system.d/hdmi-hotplug.service
[Unit]
Description=Handle HDMI hotplug
[Service]
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/handle-hdmi-hotplug
Type=oneshot
