# projects/ROCKNIX/packages/rocknix/system.d/hdmi-hotplug.path
[Unit]
Description=Monitor HDMI status

[Path]
PathModified=/run/hdmi-status.last

[Install]
WantedBy=multi-user.target
