# projects/ROCKNIX/packages/apps/rocknix-touchscreen-keyboard/system.d/touchkeyboard.service
[Unit]
Description=touchscreen keyboard service
After=input.service
Before=rocknix.target
StartLimitIntervalSec=0

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/rocknix-touchscreen-keyboard
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
