# projects/ROCKNIX/packages/sysutils/system-utils/udev.d/99-hdmi.rules
KERNEL=="card*", SUBSYSTEM=="drm", ENV{ID_PATH}=="platform-display-engine", ACTION=="add|change", RUN+="/usr/bin/hdmi_sense"
